# salt/roots/salt/supervisor/django-external-process.conf
{% set host = pillar['nginx'].hosts[sup_name] %}
{% set django = pillar['django'].djangos[sup_name] %}
{% set venv = pillar['venv'].venvs[sup_name] %}
{% set user = pillar['users'][sup_name] %}

[program:{{ sup_program }}]
directory={{ django.path }}
user={{ user.name }}
environment=DJANGO_SETTINGS_MODULE="{{ django.settings }}"
command={{ venv.path }}/bin/python {{ django.path }}{{ sup_command }} {{ django.path }}
stdout_logfile=/var/www/ombrelloni.it/log/{{ sup_program }}.log
redirect_stderr=true
autostart=true
autorestart=true

