# salt/roots/salt/supervisor/django.conf
{% set host = pillar['nginx'].hosts[sup_name] %}
{% set django = pillar['django'].djangos[sup_name] %}
{% set venv = pillar['venv'].venvs[sup_name] %}

[program:{{ sup_name }}]
directory={{ django.path }}
command={{ venv.path }}/bin/uwsgi --ini {{ django.path }}/uwsgi.ini
autostart=true
autorestart=true
stdout_logfile={{ host.root }}/log/supervisor.log
redirect_stderr=true
stopsignal=QUIT